# root/etc/services.d/adventurebot/run
#!/usr/bin/with-contenv bash
echo "Starting Bot...."
for (( ; ; )); do
	let i++
	cd /app
	python /app/main.py
done

exit $?